# configs/systemd/timers/remount_media_drive.timer.j2
[Unit]
Description=Run remount_media_drive.service every minute

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true

[Install]
WantedBy=timers.target