# configs/systemd/services/az_backup.service.j2
[Unit]
Description=AZ Vault Backup

[Service]
Type=oneshot
ExecStart=/home/{{ ansible_user_id }}/az_backup.sh