# configs/systemd/services/remount_media_drive.service.j2
[Unit]
Description=Remount media drive

[Service]
Type=oneshot
ExecStart=sh -c 'mountpoint -q /mnt/media || mount /mnt/media'